# etc/supervisor/uwsgi.conf
[program:uwsgi]
autorestart = false
autostart = true
command = /usr/local/bin/uwsgi --ini /usr/local/lib/python2.7/dist-packages/mayan/uwsgi.ini
redirect_stderr = true
user = root
